# service/ScriptService.service
[Unit]
Description=Custom Startup Script

[Service]
ExecStart=/root/FaceRec/Server/exec_app.sh

[Install]
WantedBy=default.target
