# service/ScriptService.service
[Unit]
Description=Custom Startup Script

[Service]
ExecStart=/root/FaceRecVer.2/Server/exec_app.sh

[Install]
WantedBy=default.target
